# test_monitor.timer
[Unit]
Description=Test process monitor timer

[Timer]
OnCalendar=minutely
Unit=test_monitor.service

[Install]
WantedBy=timers.target